# scripts/launchagents/com.matias.rotate-logs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matias.rotate-logs</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mmirvois/clawd/scripts/rotate-logs.sh</string>
    </array>

    <!-- Weekly log rotation for ~/clawd/logs/*.log > 5MB, keep 3 gzipped copies.
         Sundays 04:00 local (quiet window). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>4</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mmirvois/clawd/logs/rotate-logs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mmirvois/clawd/logs/rotate-logs.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
