# scripts/launchagents/com.matias.strategy-engine.paper-mtm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matias.strategy-engine.paper-mtm</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mmirvois/clawd/venv/bin/strategy-engine</string>
        <string>paper</string>
        <string>mtm</string>
    </array>

    <!-- Mark-to-market + auto-close all open paper positions.
         Runs weekdays 17:00 local (after close + fmp daily_update bypassing 7 PM lock window). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/mmirvois/clawd/logs/paper-mtm.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mmirvois/clawd/logs/paper-mtm.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
